# .config/autostart/keepassxc.desktop
[Desktop Entry]
Exec=/usr/bin/keepassxc
Name=keepassxc
Comment=a modern open-source password manager
Type=Application
